# scripts/poppit.cron
* * * * * __PROJECT_PATH__/scripts/backup.sh