# scripts/poppit.cron
0 0 * * * __PROJECT_PATH__/scripts/backup.sh
